# #dev/supervisor/local/http-server-workers.ini
[program:app-http-server-worker-0]
command=php /var/www/data/codenrock/back/app/project/console app:http-server --socket-mode 1 --worker-number 0 -q
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/supervisor/app-http-server-worker-0.log
user=apache
autostart=true
autorestart=true